# systemd/rob-report.timer
[Unit]
Description=Daily The ROB Report generation

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
